# catkin_ws/src/summit_packages_release_melodic-devel/robot_bringup/scripts/bringup-ros.service
[Unit]
Description=Launches all ROS stuff through screen
Requires=dbus.socket
After=multi-user.target

[Service]
User=robot
RemainAfterExit=yes
ExecStart=/home/robot/bringup.sh

[Install]
WantedBy=multi-user.target
